# trunk/tools/unxsvz/cron.d/unxsvz_changerootpasswd.cron
##Change hardware node passwd daily
7 3 * * * root /usr/sbin/unxsvzChangeRootPasswd.sh changepasswd > /tmp/unxsvzChangeRootPasswd.sh.log 2>&1;
