# unxsScripts/unxsvz/cron.d/unxsvz_AsteriskPeers.cron
# this file should be placed in /etc/cron.d/
#Gather active PBX extension peer IP class Cs for tIP and unxsSnort use
45 10,14,18 * * * root sleep $[($RANDOM\%300)+45]; /usr/sbin/unxsvzGatherPeerIPs.sh run > /tmp/unxsvzGatherPeerIPs.sh.log 2>&1;
