# unxsRadius/tools/crontab/proxy.root
* * * * * /usr/local/bin/unxsRadius.cgi ProcessJobQueue radius-proxy0 >> /var/log/unxsRadius.log 2>&1;
