# unxsRadius/tools/crontab/proxy.root
* * * * * /usr/local/bin/mysqlRadius2.cgi ProcessJobQueue radius-proxy0 >> /var/log/mysqlRadius2.log 2>&1;
