# unxsRadius/tools/crontab/server.root
#unxsRadius job queue processor
* * * * * /usr/local/bin/unxsRadius.cgi ProcessExtJobQueue unxsrad2 >> /var/log/mysqlradius2.log 2>&1; /usr/local/bin/unxsRadius.cgi ProcessJobQueue unxsrad2 >> /var/log/mysqlradius2.log 2>&1

