# unxsRadius/tools/crontab/server.root
#mysqlRadius2 job queue processor
* * * * * /usr/local/bin/mysqlRadius2.cgi ProcessExtJobQueue unxsrad2 >> /var/log/mysqlradius2.log 2>&1; /usr/local/bin/mysqlRadius2.cgi ProcessJobQueue unxsrad2 >> /var/log/mysqlradius2.log 2>&1

